# config/systemd/user/offlineimap-oneshot.service.d/timer.conf
# cf. https://wiki.archlinux.org/index.php/OfflineIMAP#systemd_timer

[Timer]
OnUnitInactiveSec=60m
