# config/systemd/user/borgbackup.service
[Unit]
Description=Run borgbackup.sh

[Service]
ExecStart=/home/flowfx/bin/borgbackup.sh
Type=oneshot
