# config/systemd/user/borgbackup.timer
[Unit]
Description=Run borgbackup.sh

[Timer]
OnBootSec=30m
OnUnitActiveSec=90m
AccuracySec=5m

[Install]
WantedBy=timers.target
